# crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.vyYSYI/crontab installed on Mon Jan  4 21:37:37 2021)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
* * * * * /etc/hetrixtools/hetrixtools_agent.sh >> /etc/hetrixtools/hetrixtools_cron.log 2>&1
0 0 1 * * "/etc/letsencrypt"/acme.sh --cron --home "/etc/letsencrypt" --config-home "/etc/letsencrypt/config" > /dev/null
0 0 * * * rclone move /mnt/uploadavi/ crypted:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4 --include *.avi
0 1 * * * rclone move /mnt/upload/ crypted1:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 2 * * * rclone move /mnt/upload/ crypted2:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 3 * * * rclone move /mnt/upload/ crypted3:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 4 * * * rclone move /mnt/upload/ crypted4:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 5 * * * rclone move /mnt/upload/ crypted5:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 6 * * * rclone move /mnt/upload/ crypted:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 7 * * * rclone move /mnt/upload/ crypted1:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 8 * * * rclone move /mnt/upload/ crypted2:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 9 * * * rclone move /mnt/upload/ crypted3:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 10 * * * rclone move /mnt/upload/ crypted4:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 11 * * * rclone move /mnt/upload/ crypted5:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 12 * * * rclone move /mnt/upload/ crypted:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 13 * * * rclone move /mnt/upload/ crypted1:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 14 * * * rclone move /mnt/upload/ crypted2:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 15 * * * rclone move /mnt/upload/ crypted3:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 16 * * * rclone move /mnt/upload/ crypted4:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 17 * * * rclone move /mnt/upload/ crypted5:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 18 * * * rclone move /mnt/upload/ crypted:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 19 * * * rclone move /mnt/upload/ crypted1:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 20 * * * rclone move /mnt/upload/ crypted2:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 21 * * * rclone move /mnt/upload/ crypted3:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 22 * * * rclone move /mnt/upload/ crypted4:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4
0 23 * * * rclone move /mnt/upload/ crypted5:/STUFF/ --delete-empty-src-dirs --progress --ignore-existing --transfers 4

0 1 * * * du -sh /mnt/stuff > /mnt/stuff/aktuellegroesse.txt

#0 6 * * * filebot -rename "/mnt/stuff/movies/" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
#0 18 * * * filebot -rename "/mnt/stuff/serien/" -r --db TheTVDB --conflict auto --lang de --action move --format "{n} ({y})/{'Season '+s}/{n} - {s00e00} - {t}" -no-xattr --apply artwork nfo url metadata



@reboot sleep 600 && sudo mount -t cifs -o username=Administrator,password=xx //xx/U /mnt/upload
@reboot sleep 600 && sudo mount -t cifs -o username=Administrator,password=xx //xx/L /mnt/uploadavi


#old@reboot rclone mount --allow-other --dir-cache-time 1h --drive-chunk-size 64M --log-level INFO --vfs-read-chunk-size 32M --vfs-read-chunk-size-limit off crypted:/STUFF/MOVIES /mnt/stuff/movies --vfs-cache-mode writes --fast-list --allow-non-empty
#old@reboot rclone mount --allow-other --dir-cache-time 1h --drive-chunk-size 64M --log-level INFO --vfs-read-chunk-size 32M --vfs-read-chunk-size-limit off crypted1:/STUFF/SERIEN /mnt/stuff/serien --vfs-cache-mode writes --fast-list --allow-non-empty
#old@reboot rclone mount --allow-other --dir-cache-time 1h --drive-chunk-size 64M --log-level INFO --vfs-read-chunk-size 32M --vfs-read-chunk-size-limit off crypted:/STUFF/GAMES /mnt/stuff/games --vfs-cache-mode writes --fast-list --allow-non-empty
#old@reboot rclone mount --allow-other --dir-cache-time 1h --drive-chunk-size 64M --log-level INFO --vfs-read-chunk-size 32M --vfs-read-chunk-size-limit off crypted:/STUFF/APPS /mnt/stuff/apps --vfs-cache-mode writes --fast-list --allow-non-empty
@reboot rclone mount crypted1:/STUFF/SERIEN /mnt/stuff/serien --allow-non-empty --allow-other --size-only --log-level INFO --drive-use-trash=false --timeout 1h --vfs-cache-mode writes --dir-cache-time 2m --vfs-cache-max-age 1h --vfs-read-chunk-size 128M --vfs-read-chunk-size-limit off --drive-chunk-size 128M --buffer-size 16m --fast-list

@reboot rclone mount crypted5:/STUFF/MOVIES /mnt/stuff/movies --allow-non-empty --allow-other --size-only --log-level INFO --drive-use-trash=false --timeout 1h --vfs-cache-mode writes --dir-cache-time 2m --vfs-cache-max-age 1h --vfs-read-chunk-size 128M --vfs-read-chunk-size-limit off --drive-chunk-size 128M --buffer-size 16m --fast-list

@reboot rclone mount crypted1:/STUFF/GAMES /mnt/stuff/games --allow-non-empty --allow-other --size-only --log-level INFO --drive-use-trash=false --timeout 1h --vfs-cache-mode writes --dir-cache-time 2m --vfs-cache-max-age 1h --vfs-read-chunk-size 128M --vfs-read-chunk-size-limit off --drive-chunk-size 128M --buffer-size 16m --fast-list

@reboot rclone mount crypted:/STUFF/APPS /mnt/stuff/apps --allow-non-empty --allow-other --size-only --log-level INFO --drive-use-trash=false --timeout 1h --vfs-cache-mode writes --dir-cache-time 2m --vfs-cache-max-age 1h --vfs-read-chunk-size 128M --vfs-read-chunk-size-limit off --drive-chunk-size 128M --buffer-size 16m --fast-list








#filebot-movies
0 0 1 * * filebot -rename "/mnt/stuff/movies/a" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 2 * * filebot -rename "/mnt/stuff/movies/b" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 3 * * filebot -rename "/mnt/stuff/movies/c" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 4 * * filebot -rename "/mnt/stuff/movies/d" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 5 * * filebot -rename "/mnt/stuff/movies/e" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 6 * * filebot -rename "/mnt/stuff/movies/f" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 7 * * filebot -rename "/mnt/stuff/movies/g" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 8 * * filebot -rename "/mnt/stuff/movies/h" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 9 * * filebot -rename "/mnt/stuff/movies/i" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 10 * * filebot -rename "/mnt/stuff/movies/j" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 11 * * filebot -rename "/mnt/stuff/movies/k" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 12 * * filebot -rename "/mnt/stuff/movies/l" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 13 * * filebot -rename "/mnt/stuff/movies/m" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 14 * * filebot -rename "/mnt/stuff/movies/n" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 15 * * filebot -rename "/mnt/stuff/movies/o" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 16 * * filebot -rename "/mnt/stuff/movies/p" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 17 * * filebot -rename "/mnt/stuff/movies/q" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 18 * * filebot -rename "/mnt/stuff/movies/r" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 19 * * filebot -rename "/mnt/stuff/movies/s" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 20 * * filebot -rename "/mnt/stuff/movies/t" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 21 * * filebot -rename "/mnt/stuff/movies/u" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 22 * * filebot -rename "/mnt/stuff/movies/v" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 23 * * filebot -rename "/mnt/stuff/movies/w" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 24 * * filebot -rename "/mnt/stuff/movies/x" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 25 * * filebot -rename "/mnt/stuff/movies/y" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 26 * * filebot -rename "/mnt/stuff/movies/z" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 27 * * filebot -rename "/mnt/stuff/movies/0-9" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 28 * * filebot -rename "/mnt/stuff/movies/_der" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 29 * * filebot -rename "/mnt/stuff/movies/_die" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 0 30 * * filebot -rename "/mnt/stuff/movies/_das" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 12 4 * * filebot -rename "/mnt/stuff/movies/_the" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata
0 12 5 * * filebot -rename "/mnt/stuff/movies/unsorted" -r --db TheMovieDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})" -no-xattr --apply artwork nfo url metadata




###networkstats

* * * * * /usr/bin/vnstati -s -i enp35s0+br0+vnet0 -o /mnt/stuff/summary.png
* * * * * /usr/bin/vnstati -h -i enp35s0+br0+vnet0 -o /mnt/stuff/hour.png
* * * * * /usr/bin/vnstati -m -i enp35s0+br0+vnet0 -o /mnt/stuff/monthly.png
53 10 * * * /etc/webmin/package-updates/update.pl









#filebot-serien
0 12 1 * * filebot -rename "/mnt/stuff/serien/a" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 2 * * filebot -rename "/mnt/stuff/serien/b" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 3 * * filebot -rename "/mnt/stuff/serien/c" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 4 * * filebot -rename "/mnt/stuff/serien/d" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 5 * * filebot -rename "/mnt/stuff/serien/e" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 6 * * filebot -rename "/mnt/stuff/serien/f" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 7 * * filebot -rename "/mnt/stuff/serien/g" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 8 * * filebot -rename "/mnt/stuff/serien/h" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 9 * * filebot -rename "/mnt/stuff/serien/i" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 10 * * filebot -rename "/mnt/stuff/serien/j" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 11 * * filebot -rename "/mnt/stuff/serien/k" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 12 * * filebot -rename "/mnt/stuff/serien/l" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 13 * * filebot -rename "/mnt/stuff/serien/m" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 14 * * filebot -rename "/mnt/stuff/serien/n" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 15 * * filebot -rename "/mnt/stuff/serien/o" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 16 * * filebot -rename "/mnt/stuff/serien/p" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 17 * * filebot -rename "/mnt/stuff/serien/q" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 18 * * filebot -rename "/mnt/stuff/serien/r" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 19 * * filebot -rename "/mnt/stuff/serien/s" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 20 * * filebot -rename "/mnt/stuff/serien/t" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 21 * * filebot -rename "/mnt/stuff/serien/u" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 22 * * filebot -rename "/mnt/stuff/serien/v" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 23 * * filebot -rename "/mnt/stuff/serien/w" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 24 * * filebot -rename "/mnt/stuff/serien/x" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 25 * * filebot -rename "/mnt/stuff/serien/y" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 26 * * filebot -rename "/mnt/stuff/serien/z" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
0 12 27 * * filebot -rename "/mnt/stuff/serien/0-9" -r --db TheTVDB --conflict auto --lang de --action move --format "{n.ascii()} ({y})/{'Staffel '+s}/{n.ascii()} - {s00e00} - {t.ascii()}" -no-xattr --apply artwork nfo url metadata
#{15397b2b-cb23-467c-986f-1789505e06c2}# Veeam backup job, do not modify this line
0 6 * * * veeamconfig job start --id {15397b2b-cb23-467c-986f-1789505e06c2} --retriable --scheduled --highpriority > /dev/null
